# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:python3]
command=/bin/bash -c "python3 WLCG-site-snmp-http.py"
autostart=true
autorestart=true
